# systemd_units/djing_backup.timer
[Unit]
Description=Run backup periodically

[Timer]
OnCalendar=*-*-* 8,12,14,16,19,23:15:0
Unit=djing_backup.service

[Install]
WantedBy=timers.target

